# ping_websites.service
[Unit]
Description=Regularly ping my websites

[Service]
ExecStart=/path/of/your/ping_websites.sh

[Install]
WantedBy=multi-user.target